# override/airootfs/etc/skel/.config/autostart/change-wallpaper.desktop
[Desktop Entry]
Type=Application
Icon=folder
Name=Set Andromeda wallpaper
Exec=/usr/local/lib/andromeda/change-wallpaper.sh
Terminal=false
Hidden=false
